# .config/sv/wg/finish
#!/bin/sh
set -e

# name=$(pwd | awk -F- '{ print $2 }')
name=$(cat NAME)

echo $name

wg-quick down "$name"
