# .config/sv/mpd/run
#!/bin/sh

MPDCONF=$HOME/.config/mpd/mpd.conf
exec mpd --no-daemon $MPDCONF
